# supervisor/nova.service
[Unit]
Description=Nova Supervisor Bot
After=network.target atlas.service
Wants=atlas.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nova
ExecStart=/opt/nova/target/release/nova /opt/nova/data/prod/nova.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nova

# Allow full system access (intentional — supervisor needs it)
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
